# calibre-web.service
[Unit]
Description=UNCaGED Dashboard - Calibre Web Management
After=network.target

[Service]
Type=simple
User=wolf
Group=wolf
WorkingDirectory=/home/wolf/CODE/Python/Ebook/KOBO/calibre-web-management
EnvironmentFile=/home/wolf/CODE/Python/Ebook/KOBO/calibre-web-management/.env
ExecStart=/home/wolf/CODE/Python/Ebook/KOBO/calibre-web-management/venv/bin/python3 server_prod.py
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=calibre-web

[Install]
WantedBy=multi-user.target
